# deploy/systemd/planner-task-reminders.service
[Unit]
Description=Planner Task Reminders Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=planner
Group=planner
WorkingDirectory=/opt/planner
EnvironmentFile=/etc/planner/planner.env
ExecStart=/usr/bin/npm run -w apps/api task-reminders:worker
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
